# config/chatbot.service
[Unit]
Description=Rasa api.
After=network.target

[Service]
User=$USER
Group=www-data
WorkingDirectory=/srv/backend/
Environment="PATH=/srv/venvs/backend/bin/"
ExecStart=/srv/venvs/backend/bin/rasa run --enable-api

[Install]
WantedBy=multi-user.target
